# adaptive-firewall.service
[Unit]
Description=Adaptive DCF Firewall (NFQUEUE)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/adcfw/adaptive_firewall.py
WorkingDirectory=/opt/adcfw
Restart=on-failure
RestartSec=2
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN
# Ensure control socket dir exists
ExecStartPre=/bin/mkdir -p /var/run
User=root

[Install]
WantedBy=multi-user.target
